# scripts/motion-detector.service
[Unit]
Description=Raspberry Pi Motion Detection System
Documentation=https://github.com/dylanpatriarchi/raspberry-pi-motion-detection
After=network.target

[Service]
Type=simple

# Adjust User/WorkingDirectory to match your deployment.
User=pi
WorkingDirectory=/home/pi/raspberry-pi-motion-detection

# Run headless (no preview window) as a background service.
ExecStart=/usr/bin/python3 main.py --no-preview

Restart=on-failure
RestartSec=5

# Basic hardening.
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
